# cron/uninsubria_courses.cron
# Daily at 07:00: sync courses -> notebooks -> folder files -> NotebookLM sources
0 7 * * * /Users/riccardosacchetto/Desktop/BIG\ DATA/scripts/sync_uninsubria_notebooks.sh
